# deploy/cloudflared.service
[Unit]
Description=Cloudflare Tunnel for portfolio-ai
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/cloudflared tunnel run \
  --url http://127.0.0.1:8000
Restart=on-failure
RestartSec=5
Environment=HOME=/root

[Install]
WantedBy=multi-user.target
